# pic.conf
[program:pic]
directory=/root/pic
command=/usr/bin/python3.6 /root/pic/app.py
